# maro/cli/grass/lib/agents/systemd/maro-node-agent.service
[Unit]
Description=Maro Node Agent
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
WorkingDirectory={home_path}/.maro/lib/grass
ExecStart=/usr/bin/python3 -m agents.node_agent

[Install]
WantedBy=default.target
